# swarmchat-api-prod.service
[Unit]
Description=SwarmChat FastAPI Application (Production)
After=network.target
Wants=network-online.target

[Service]
User=peter
Group=users
WorkingDirectory=/home/peter/convo
Environment=LANG=en_US.UTF-8
Environment=LC_ALL=en_US.UTF-8
EnvironmentFile=/home/peter/convo/.env
ExecStart=/home/peter/.virtualenvs/swarm/bin/uvicorn swarm_chat:app --host 127.0.0.1 --port 8000 --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
TimeoutStartSec=30
TimeoutStopSec=5
Type=simple
KillMode=mixed

[Install]
WantedBy=multi-user.target
